# App/mavtcp_working.service
[Unit]
Description=Pixhawk MAVLink TCP Bridge (Working)
After=network.target
StartLimitInterval=0

[Service]
Type=simple
# ÇALIŞAN SOCAT KOMUTU - Serial RAW mode
ExecStart=/usr/bin/socat -d -d TCP-LISTEN:5777,bind=0.0.0.0,reuseaddr,fork /dev/ttyACM0,b115200,raw,echo=0,clocal=1,cs8

Restart=always
RestartSec=2
User=root
StandardOutput=journal
StandardError=journal

# Port hazır olana kadar bekle
ExecStartPre=/bin/bash -c 'until [ -c /dev/ttyACM0 ]; do sleep 1; done'

[Install]
WantedBy=multi-user.target 