# App/mavtcp_fixed.service
[Unit]
Description=Pixhawk MAVLink TCP Bridge (Fixed)
After=network.target
StartLimitInterval=0

[Service]
Type=simple
# DÜZELTİLMİŞ SOCAT KOMUTU - raw terminal ayarları kaldırıldı
ExecStart=/usr/bin/socat -d -d TCP-LISTEN:5777,bind=0.0.0.0,reuseaddr,fork OPEN:/dev/ttyACM0,nonblock
# Alternatif 1: Eğer yukarısı çalışmazsa
# ExecStart=/usr/bin/socat TCP-LISTEN:5777,bind=0.0.0.0,reuseaddr,fork /dev/ttyACM0,b115200,cs8,parenb=0,cstopb=0,clocal=1,raw,echo=0
# Alternatif 2: Eğer ttyUSB0 ise
# ExecStart=/usr/bin/socat TCP-LISTEN:5777,bind=0.0.0.0,reuseaddr,fork /dev/ttyUSB0,b115200,cs8,parenb=0,cstopb=0,clocal=1,raw,echo=0

Restart=always
RestartSec=2
User=root
StandardOutput=journal
StandardError=journal

# Port hazır olana kadar bekle
ExecStartPre=/bin/bash -c 'until [ -e /dev/ttyACM0 ]; do sleep 1; done'

[Install]
WantedBy=multi-user.target 